# systemd/tuya-monitoring.service
[Unit]
Description=Tuya Smart Plug Monitoring (Local LAN → Graphite)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=nickc
WorkingDirectory=/home/nickc/code/electricity_monitoring
ExecStart=/usr/bin/python3 /home/nickc/code/electricity_monitoring/tuya_local_to_graphite.py
Restart=always
RestartSec=10
StandardOutput=append:/home/nickc/electricity_tuya_local.log
StandardError=append:/home/nickc/electricity_tuya_local.log

[Install]
WantedBy=multi-user.target
